# systemd/sunset-cam-boot.service
[Unit]
Description=Sunset Cam boot dispatcher (decides SETUP vs ONLINE)
# Must run AFTER NetworkManager has finished its startup connection attempt — not
# merely after it has *started*. `After=NetworkManager.service` only means NM's
# process launched; it has NOT yet loaded saved profiles or attempted autoconnect,
# so nmcli returns an empty list and has_wifi_credentials() wrongly reports "no
# creds" → a false SETUP decision (observed on a slow Pi Zero W: dispatcher ran in
# 5s, hit SETUP, then NM joined home WiFi → AP-vs-home flap).
# network-online.target (reached when NetworkManager-wait-online completes) blocks
# until NM is settled — connected, or timed out — so by the time we run:
#   1. nmcli lists the real saved WiFi connections (has_wifi_credentials), and
#   2. a good home network is already associated (is_online true on the first poll).
After=network-online.target NetworkManager.service
Wants=network-online.target
# Deliberately NO `Before=` the units we start. boot.main() starts the chosen unit
# imperatively (a blocking `systemctl start`). If this oneshot were ordered Before=
# them, systemd would refuse to start an ordered-after unit from within our own
# still-running ExecStart and silently drop the job — the device boots, joins WiFi,
# and sits dark. The dispatcher owns sequencing; no ordering edge is needed here.

[Service]
Type=oneshot
# dispatch_boot() in boot.main():
#   1. Checks whether saved home WiFi creds exist (has_wifi_credentials).
#   2. If no creds → starts sunset-cam-setup.service (captive-portal SETUP).
#   3. If creds exist → polls is_online() up to ~60 s for NM to associate.
#      Success → starts sunset-cam-supervisor.service (ONLINE).
#      Timeout (bad password / network down) → falls back to
#      sunset-cam-setup.service (SETUP) so the customer can re-enter creds.
ExecStart=/opt/sunset-cam/.venv/bin/python -m sunset_cam.boot
RemainAfterExit=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
