# systemd/sunset-cam-boot.service
[Unit]
Description=Sunset Cam boot dispatcher (decides SETUP vs ONLINE)
# ⚠️  HARDWARE-GATED ordering notes:
#   - After=basic.target: correct ordering for most Pi OS Lite images.
#     Adjust to network-pre.target if you want it to run after wlan0 is up
#     (so wpa_supplicant has already attempted association before the check).
#   - Before=sunset-cam-supervisor.service sunset-cam-setup.service ensures
#     this runs first and lets dispatch_boot start the right target.
After=basic.target
Before=sunset-cam-supervisor.service sunset-cam-setup.service

[Service]
Type=oneshot
# dispatch_boot() in boot.main() checks /etc/wpa_supplicant/wpa_supplicant.conf
# for a network={ block; if absent → starts sunset-cam-setup.service (SETUP);
# if present → starts sunset-cam-supervisor.service (ONLINE).
ExecStart=/opt/sunset-cam/.venv/bin/python -m sunset_cam.boot
RemainAfterExit=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
