# systemd/sunset-cam-boot.service
[Unit]
Description=Sunset Cam boot dispatcher (decides SETUP vs ONLINE)
# NetworkManager must be up before this runs so that:
#   1. nmcli can list saved WiFi connections (has_wifi_credentials).
#   2. nmcli can check active connections (is_online) while we wait for the
#      join to complete.
# Wants= pulls NM in if it is not already in the boot graph.
After=NetworkManager.service
Wants=NetworkManager.service
Before=sunset-cam-supervisor.service sunset-cam-setup.service

[Service]
Type=oneshot
# dispatch_boot() in boot.main():
#   1. Checks whether saved home WiFi creds exist (has_wifi_credentials).
#   2. If no creds → starts sunset-cam-setup.service (captive-portal SETUP).
#   3. If creds exist → polls is_online() up to ~60 s for NM to associate.
#      Success → starts sunset-cam-supervisor.service (ONLINE).
#      Timeout (bad password / network down) → falls back to
#      sunset-cam-setup.service (SETUP) so the customer can re-enter creds.
ExecStart=/opt/sunset-cam/.venv/bin/python -m sunset_cam.boot
RemainAfterExit=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
